# piwriter.service
[Unit]
Description=PiWriter display driver
After=multi-user.target

[Service]
WorkingDirectory=/home/pi/piwriter
Type=idle
ExecStart=/usr/bin/python3 /home/pi/piwriter/main.py run
ExecStop=/usr/bin/python3 /home/pi/piwriter/main.py sleep
# Restart=always

[Install]
WantedBy=multi-user.target
